# systemd/opencode-serve.service
[Unit]
Description=OpenCode Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/fen/.npm-global/bin/opencode serve --port 4096 --print-logs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
